# deploy/illumio-ops.service
[Unit]
Description=Illumio PCE Ops (offline bundle)
Documentation=https://github.com/your-org/illumio_ops
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=illumio_ops
Group=illumio_ops
WorkingDirectory=/opt/illumio_ops
ExecStart=/opt/illumio_ops/python/bin/python3 /opt/illumio_ops/illumio_ops.py --monitor-gui --interval 10
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=illumio-ops

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/illumio_ops

[Install]
WantedBy=multi-user.target
