# deploy/illumio-ops.service
[Unit]
Description=Illumio PCE Ops
Documentation=https://github.com/your-org/illumio_ops
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=illumio_ops
Group=illumio_ops
WorkingDirectory=/opt/illumio_ops
# Ubuntu/Debian: PEP 668 blocks system-wide pip; use a venv.
#   python3 -m venv /opt/illumio_ops/venv
#   /opt/illumio_ops/venv/bin/pip install flask pandas pyyaml
ExecStart=/opt/illumio_ops/venv/bin/python /opt/illumio_ops/illumio_ops.py --monitor --interval 10
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=illumio-ops

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/illumio_ops

[Install]
WantedBy=multi-user.target
